# altserver-iphone-backup/s6-rc.d/dbus/run
#!/command/execlineb -P
pipeline -w { sed --unbuffered "s/^/[dbus] /" }
fdmove -c 2 1
dbus-daemon --config-file=/usr/share/dbus-1/system.conf --print-address --nofork