# altserver-iphone-backup/s6-rc.d/anisette/run
#!/command/execlineb -P
pipeline -w { sed --unbuffered "s/^/[anisette] /" }
fdmove -c 2 1
with-contenv
cd /opt
/opt/anisette_server -r -p 6969