# altserver-iphone-backup/s6-rc.d/usbmuxd/run
#!/command/execlineb -P
pipeline -w { sed --unbuffered "s/^/[usbmuxd] /" }
fdmove -c 2 1
usbmuxd -f 