# altserver-iphone-backup/s6-rc.d/netmuxd/run
#!/command/execlineb -P
pipeline -w { sed --unbuffered "s/^/[netmuxd] /" }
fdmove -c 2 1
/opt/netmuxd --disable-unix --host 127.0.0.1