# altserver-iphone-backup/s6-rc.d/avahi/run
#!/command/execlineb -P
pipeline -w { sed --unbuffered "s/^/[avahi] /" }
fdmove -c 2 1
avahi-daemon